# systemd/arlo-viewer.service
[Unit]
Description=Arlo Recording Viewer
After=network.target

[Service]
Type=simple
# Change to your username
User=YOUR_USERNAME
# Change to where you installed arlo-viewer
WorkingDirectory=/home/YOUR_USERNAME/arlo-viewer
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
