# systemd/security-bore-tunnel.service
[Unit]
Description=Security Viewer Bore Tunnel (arlo-viewer → remote)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/home/YOUR_USERNAME
ExecStart=/usr/local/bin/bore local 3003 --to YOUR_BORE_SERVER --port YOUR_BORE_VIEWER_PORT
Restart=always
RestartSec=10
StandardOutput=append:/var/log/security-bore-tunnel.log
StandardError=append:/var/log/security-bore-tunnel.log
TimeoutStartSec=30
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
